# roles/gunicorn/templates/gunicorn.socket.j2
# -*- mode: systemd -*-
# vi: set ft=systemd syntax=systemd :
[Unit]
Description=gunicorn socket

[Socket]
ListenStream={{ gunicorn['sock_file'] }}

[Install]
WantedBy=sockets.target
